# twoflower/api/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=python main.py
