# twoflower/redis/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server /etc/redis/redis.conf --protected-mode no

[program:app]
command=python main.py
